# etc/cron.d/froxlor-system
13 0,4,8,12,16,20 * * * root /usr/bin/nice -n 5 '__INSTALLDIR__/update.php' >'__INSTALLDIR__/update.log' 2>&1
13 5 * * * root /usr/bin/nice -n 5 '__INSTALLDIR__/backup.sh' >'__INSTALLDIR__/backup.log' 2>&1
