# etc/cron.d/froxlor-system
13 3 * * * root /usr/bin/nice -n 5 ${INSTALLDIR}/backup.sh > /dev/null
13 6 * * * root /usr/bin/nice -n 5 ${INSTALLDIR}/update.sh > /dev/null
