# etc/cron.d/froxlor-system
13 */4 * * * root /usr/bin/nice -n 5 '${INSTALLDIR}/update.php' >'${INSTALLDIR}/update.log' 2>&1
13 3 * * * root /usr/bin/nice -n 5 '${INSTALLDIR}/backup.sh' >'${INSTALLDIR}/backup.log' 2>&1
